# root/etc/services.d/plex_sync/run
#!/usr/bin/with-contenv bash

while :
do
  flock -n "$SYNC_STATUS_FILE" /usr/local/bin/sync_and_log "$HOME" "$HOME_PERSISTENT"
  sleep 60s
done
